# {{cookiecutter.project_slug}}/etc/celerybeat.conf
{% raw -%}
[program:{{ctx.project_slug}}-celerybeat]
directory={{ctx.django_dir}}
command={{ctx.virtualenv_dir}}/bin/celery beat -A {{ctx.project_slug}} --schedule {{ctx.run_dir}}/{{ctx.project_slug}}-beat.db -l info
user={{ctx.host_user}}
autostart=true
autorestart=true
stdout_logfile={{ctx.log_dir}}/{{ctx.project_slug}}_celerybeat.log
redirect_stderr=true
{%- endraw %}
