# {{cookiecutter.project_slug}}/etc/celery.conf
{% raw -%}
[program:{{ctx.project_slug}}-celery]
directory={{ctx.django_dir}}
command={{ctx.virtualenv_dir}}/bin/celery worker -A {{ctx.project_slug}} -E -l info --concurrency=4
user={{ctx.host_user}}
autostart=true
autorestart=true
stdout_logfile={{ctx.log_dir}}/{{ctx.project_slug}}_celery.log
redirect_stderr=true
{%- endraw %}
