# {{cookiecutter.project_slug}}/etc/uwsgi.conf
{% raw -%}
[program:{{ctx.project_slug}}-uwsgi]
directory={{ctx.install_dir}}
command={{ctx.virtualenv_dir}}/bin/uwsgi --ini {{ctx.django_dir}}/uwsgi.ini
user={{ctx.host_user}}
autostart=true
autorestart=true
stdout_logfile={{ctx.log_dir}}/{{ctx.project_slug}}_uwsgi.log
redirect_stderr=true
{%- endraw %}
